# deploy/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0777

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=error
nodaemon=false
minfds=1024
minprocs=200
user=vt
childlogdir=/tmp/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:vt]
command=/home/vt/vt-env/bin/gunicorn deploy.wsgi:application --workers 2 --worker-class gevent -c /home/vt/vt-app/config/gunicorn_conf.py --timeout 360
numprocs=1
process_name=%(program_name)s
directory=/home/vt/vt-app
user=vt
stderr_logfile = /home/vt/vt-stderr.log
stdout_logfile = /home/vt/vt-stdout.log
